# Build/com.grahamgilbert.crypt.launcher.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.grahamgilbert.crypt.launcher</string>
	<key>ProgramArguments</key>
	<array>
		<string>/usr/local/crypt/Crypt.app/Contents/MacOS/Crypt</string>
	</array>
	<key>WatchPaths</key>
	<array>
		<string>/usr/local/crypt/watch</string>
	</array>
</dict>
</plist>
