# Build/com.aapps.crypt.launcher.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
<<<<<<< HEAD
	<string>com.aapps.crypt.launcher</string>
=======
	<string>com.aapps.crypt-launcher</string>
>>>>>>> 55c96dc0fab9d6f1d6bf5f0cfcec17d73296e254
	<key>ProgramArguments</key>
	<array>
		<string>/usr/local/crypt/Crypt.app/Contents/MacOS/Crypt</string>
	</array>
	<key>WatchPaths</key>
	<array>
		<string>/tmp/crypt-launcher/</string>
	</array>
</dict>
</plist>
